# deployment/templates/systemd/sage-daemon.service
[Unit]
Description=SAGE JobManager Daemon Service
Documentation=https://github.com/intellistream/SAGE
After=sage-ray.service
Wants=sage-ray.service

[Service]
Type=simple
User=sage
Group=sage
WorkingDirectory=/opt/sage
Environment=DAEMON_HOST=127.0.0.1
Environment=DAEMON_PORT=19001
Environment=ACTOR_NAME=sage_global_jobmanager
Environment=NAMESPACE=sage_system
ExecStart=/opt/sage/deployment/scripts/daemon_manager.sh start
ExecStop=/opt/sage/deployment/scripts/daemon_manager.sh stop
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
