# deployment/templates/systemd/sage-ray.service
# SAGE systemd service templates
# 系统服务模板

[Unit]
Description=SAGE Ray Head Node Service
Documentation=https://github.com/intellistream/SAGE
After=network.target
Wants=network.target

[Service]
Type=forking
User=sage
Group=sage
WorkingDirectory=/opt/sage
Environment=RAY_HEAD_PORT=10001
Environment=RAY_DASHBOARD_PORT=8265
Environment=RAY_TEMP_DIR=/var/lib/ray_shared
ExecStart=/opt/sage/deployment/scripts/ray_manager.sh start
ExecStop=/opt/sage/deployment/scripts/ray_manager.sh stop
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
